# launchd.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-assistant.watcher</string>

    <key>Program</key>
    <string>{{REPO_DIR}}/.venv/bin/python</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{REPO_DIR}}/.venv/bin/python</string>
        <string>{{REPO_DIR}}/watcher.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{REPO_DIR}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{BOOTSTRAP_PATH}}</string>
    </dict>

    <key>StandardOutPath</key>
    <string>{{REPO_DIR}}/logs/watcher.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{REPO_DIR}}/logs/watcher.stderr.log</string>
</dict>
</plist>
